# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
; /static dilayani WhiteNoise di Django (lihat settings); jangan static-map uwsgi agar tidak 404 saat disk kosong / volume menimpa.
; Increase request/header buffers to handle large Cookie headers (e.g., passkey session flow)
; and prevent "invalid request block size: ... (max 4096)" that surfaces as 502 at proxy.
command=uwsgi --http :8001 --module sso_service.wsgi:application --master --processes 4 --threads 2 --buffer-size 65535 --http-buffer-size 65535
directory=/usr/src/app
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi.log
stderr_logfile=/var/log/uwsgi.err

[program:qcluster]
command=python manage.py qcluster
directory=/usr/src/app
autostart=true
autorestart=true
stdout_logfile=/var/log/qcluster.log
stderr_logfile=/var/log/qcluster.err
